# sysadmin/systemd/shopware_sitemap.timer
[Unit]
Description=Shopware: Generate sitemap

[Timer]
OnCalendar=daily
Unit=shopware_sitemap.service

[Install]
WantedBy=timers.target